# powertune.service
[Unit]
Description=Power tuning

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/powertune.sh

[Install]
WantedBy=multi-user.target

